# files/etc/service/docker-gen-letsencrypt-hosts/run
#!/usr/bin/env bash

exec docker-gen -watch /app/le-hosts.tmpl /config/letsencrypt/hosts.txt
